# rootfs/etc/services.d/cron/run
#!/usr/bin/with-contenv sh

exec crond -f
